# machine/desktop/root/etc/supervisor.d/echo.conf
[program:echo_example]
directory=/
command=echo Hello, World!
